# infra/crontab
# Crontab biz-hunter.
# Cron de Linux: el container `biz-hunter-cron` dispara los scripts.
# Vars de entorno se cargan desde /opt/biz-hunter/.env via env-file en docker-compose.
# Output va a /var/log/biz-hunter.log → linkado a stdout, visible con `docker logs`.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Scout: cada 6h
0 */6 * * * root cd /opt/biz-hunter && python scripts/scout.py >> /var/log/biz-hunter.log 2>&1

# Validator: 4 veces al día (cada 6h, offset 30min)
30 */6 * * * root cd /opt/biz-hunter && python scripts/validator.py >> /var/log/biz-hunter.log 2>&1

# Notify validated → Telegram con botones aprobar/rechazar
# Cada 4h, offset 45min para no chocar con validator. Sólo manda las
# validadas con notified_at IS NULL (anti-duplicados natural).
45 */4 * * * root cd /opt/biz-hunter && python scripts/notify_validated.py >> /var/log/biz-hunter.log 2>&1

# Builder: cada 30min, procesa la cola de approved (FIFO, max 5/run)
*/30 * * * * root cd /opt/biz-hunter && python scripts/builder.py >> /var/log/biz-hunter.log 2>&1

# Tester: 3 veces al día (cada 8h)
0 */8 * * * root cd /opt/biz-hunter && python scripts/tester.py >> /var/log/biz-hunter.log 2>&1

# Reporter digest diario: 8am
0 8 * * * root cd /opt/biz-hunter && python scripts/reporter.py >> /var/log/biz-hunter.log 2>&1
